# systemd_units/CcspPandMSsp.service
[Unit]
Description=CcspPandMSsp service

After=dbus.service CcspCrSsp.service logagent.service PsmSsp.service

Requires=CcspCrSsp.service
Wants=dbus.service logagent.service PsmSsp.service

[Service]
Type=forking
WorkingDirectory=/usr/ccsp/pam
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
ExecStart=/usr/bin/CcspPandMSsp -subsys $Subsys
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
