# systemd_units/ccspwifiagent.service
[Unit]
Description=Cosa service
After=CcspPandMSsp.service

[Service]
Type=simple
WorkingDirectory=/usr/ccsp
ExecStart=/bin/sh -c '/usr/ccsp/wifi/cosa_start_wifiagent.sh &' 
RemainAfterExit=yes

StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
