# systemd_units/CcspMtaAgentSsp.service
[Unit]
Description=CcspMtaAgentSsp service

After=CcspPandMSsp.service notifyComp.service

[Service]
Type=forking
PIDFile=/var/tmp/CcspMtaAgent.pid
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
WorkingDirectory=/usr/ccsp/mta
ExecStart=/usr/bin/CcspMtaAgentSsp -subsys $Subsys
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=ccspCr-started.target
