# systemd_units/CcspCMAgentSsp.service
[Unit]
Description=CcspCMAgentSsp service

After=CcspPandMSsp.service CcspMtaAgentSsp.service

[Service]
Type=forking
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
WorkingDirectory=/usr/ccsp/cm
ExecStart=/usr/bin/CcspCMAgentSsp -subsys $Subsys
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=ccspCr-started.target
