# systemd_units/CcspCrSsp.service
[Unit]
Description=CcspCrSsp service

[Service]
Type=forking
WorkingDirectory=/usr/ccsp
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
TimeoutStartSec=120
#RDK-7248: Temp Workaround 
ExecStartPre=/bin/sh -c '(sleep 60 && /usr/ccsp/ccspSysConfigEarly.sh)'
ExecStart=/usr/bin/CcspCrSsp -subsys $Subsys
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
