# systemd_units/CcspTr069PaSsp.service
[Unit]
Description=CcspTr069PaSsp service

After=harvester.service

[Service]
Type=forking
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
EnvironmentFile=/etc/device.properties
WorkingDirectory=/usr/ccsp/tr069pa
ExecStart=/usr/bin/CcspTr069PaSsp -subsys $Subsys
ExecStop=/bin/sh -c 'echo "`date`: Stopping/Restarting CcspTr069PaSsp" >> ${PROCESS_RESTART_LOG}'
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=wifi-initialized.target
