# systemd_units/CcspWebUI.service
[Unit]
Description=CcspWebUI service
After=cosastart-rem.service

[Service]
Type=forking
WorkingDirectory=/usr/www
ExecStart=/bin/sh -c '/etc/webgui.sh'
RemainAfterExit=yes

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
