# systemd_units/CcspWebUI.service
[Unit]
Description=CcspWebUI service
After=cosastart.service cosastart-rem.service
Requires=cosastart.service cosastart-rem.service

[Service]
Type=forking
ExecStart=/bin/sh -c '/etc/webgui.sh &'
RemainAfterExit=yes

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
