# systemd_units/CcspPandMSsp.service
[Unit]
Description=CcspPandMSsp service
After=gwprovapp.service CcspCrSsp.service PsmSsp.service

[Service]
Type=forking
EnvironmentFile=/etc/utopia/service.d/log_env_var.sh
Environment="Subsys=eRT."
WorkingDirectory=/usr/ccsp/pam
ExecStart=/usr/bin/CcspPandMSsp -subsys $Subsys
ExecStop=/bin/kill -15 $MAINPID
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
