# systemd_units/CcspWecbController.service
[Unit]
Description=CcspWecbController service

After=dbusOverride.service CcspCrSsp.service logagent.service PsmSsp.service CcspPandMSsp.service
Requires=CcspCrSsp.service

[Service]
Type=forking
WorkingDirectory=/usr/ccsp/wecb
Environment="Subsys=eRT."
ExecStart=/usr/bin/CcspWecbController -subsys $Subsys
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
