# systemd_units/ccspwifiagent.service
[Unit]
Description=Cosa service
After=dbusOverride.service CcspCrSsp.service logagent.service PsmSsp.service CcspPandMSsp.service \
CcspWecbController.service
Before=CcspLMLite.service
Requires=CcspCrSsp.service CcspPandMSsp.service

[Service]
Type=oneshot
WorkingDirectory=/usr/ccsp
ExecStartPre=/bin/sh -c 'sleep 30'
ExecStart=/bin/sh -c '/usr/ccsp/wifi/cosa_start_wifiagent.sh &' 
RemainAfterExit=yes

StandardOutput=syslog

[Install]
WantedBy=multi-user.target

