# systemd_units/CcspCrSsp.service
[Unit]
Description=CcspCrSsp service

After=docsis.conf dbusOverride.service CcspSysConfig.service

[Service]
Type=forking
WorkingDirectory=/usr/ccsp
Environment="Subsys=eRT."
ExecStartPre=/bin/sh -c 'sleep 15'
ExecStart=/usr/bin/CcspCrSsp -subsys $Subsys
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
