# systemd_units/CcspCrSsp.service
[Unit]
Description=CcspCrSsp service

[Service]
Type=notify
NotifyAccess=main
PIDFile=/var/tmp/CcspCrSsp.pid
WorkingDirectory=/usr/ccsp
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
EnvironmentFile=/etc/device.properties
ExecStartPre=/bin/sh -c '(/usr/ccsp/ccspSysConfigEarly.sh)'
ExecStart=/usr/bin/CcspCrSsp -subsys $Subsys
ExecStop=/bin/sh -c 'echo 0 >> /tmp/CcspCrSsp_Restarted'
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
