# systemd_units/crResetDetect.path
[Unit]
Description=CcspCrSsp Restart Occurred Path

[Path]
PathExists=/tmp/CcspCrSsp_Restarted
Unit=crResetDetect.service

[Install]
WantedBy=multi-user.target
