# systemd_units/docsisinitialized.path
[Unit]
Description=Docsis Initialized Path

[Path]
PathExists=/var/tmp/eSafeCfgFile

[Install]
WantedBy=multi-user.target

