# systemd_units/gwprovapp.service
[Unit]
Description=Gwprovapp service

#Note: Requires, Wants, and After are device specific. see gwprovapp.conf

[Service]
EnvironmentFile=-/etc/device.properties
ExecStart=/usr/bin/gw_prov_utopia $PEER_COM_IP
Restart=always
StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
