# systemd_units/gwprovapp.service
[Unit]
Description=Gwprovapp service

#Note: Requires, Wants, and After are device specific. see gwprovapp.conf

[Service]
#Note: Don't Uncomment Below Line: Each Device needs a conf file to launch the binary
#ExecStart=/usr/bin/gw_prov_utopia $IP
Restart=always
StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
