# systemd_units/harvester.service
[Unit]
Description=harvester service

[Service]
Type=forking
WorkingDirectory=/usr/ccsp/harvester
ExecStart=/usr/bin/harvester
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=wifi-initialized.target
