# systemd_units/logagent.service
[Unit]
Description=CcspLogent service

[Service]
Type=notify
NotifyAccess=main
PIDFile=/var/tmp/log_agent.pid
WorkingDirectory=/usr/ccsp/logagent
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
ExecStart=/usr/bin/log_agent -subsys $Subsys
Restart=always

StandardOutput=syslog

[Install]
WantedBy=multi-user.target
