# systemd_units/pam-initialized.target
[Unit]
Description=PAM Initialization Complete
Documentation=man:systemd.special(7)
