# systemd_units/paminitialized.service
[Unit]
Description=Pam Initialized service


[Service]
Type=simple
ExecStart=/bin/echo "Pam is initialized"
