# systemd_units/paminitialized.path
[Unit]
Description=Pam Initialized Path

[Path]
PathExists=/tmp/pam_initialized

[Install]
WantedBy=multi-user.target

