# systemd_units/PsmSsp.service
[Unit]
Description=PsmSsp service

After=CcspCrSsp.service logagent.service

[Service]
Type=notify
NotifyAccess=main
PIDFile=/var/tmp/PsmSsp.pid
WorkingDirectory=/usr/ccsp
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
ExecStart=/usr/bin/PsmSsp -subsys $Subsys
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=multi-user.target
