# systemd_units/PsmSsp.service
[Unit]
Description=PsmSsp service

After=logagent.service

[Service]
Type=forking
PIDFile=/var/tmp/PsmSsp.pid
WorkingDirectory=/usr/ccsp
Environment="Subsys=eRT."
Environment="LOG4C_RCPATH=/etc"
ExecStart=/usr/bin/PsmSsp -subsys $Subsys
Restart=always

StandardOutput=syslog+console

[Install]
WantedBy=ccspCr-started.target
