# systemd_units/wifiinitialized.service
[Unit]
Description=Wifi Initialized Service


[Service]
Type=simple
ExecStart=/bin/echo "Wifi is initialized"
