# systemd_units/wifiinitialized.path
[Unit]
Description=Wifi Initialized Path

[Path]
PathExists=/tmp/wifi_initialized

[Install]
WantedBy=multi-user.target

